# roles/usrv-delay/files/delay.service
[Unit]
Description=Delay user session

[Service]
Type=oneshot
ExecStart=/bin/sh -c "sleep 6"

[Install]
WantedBy=default.target
