# roles/usrv-x11vnc/files/x11vnc.service
[Unit]
Description=X11vnc desktop sharing

[Service]
ExecStart=/usr/bin/x11vnc
Restart=on-abnormal

[Install]
WantedBy=delay.target
